# piThoughts.cron
MAILTO=""
*/10	*	*	*	*	python3 /home/pi/thoughts.py --colour red >> /home/pi/thoughts.log
