# piThoughts.cron
MAILTO=""
*/10	*	*	*	*	root	python3 /home/pi/thoughts.py --colour red
